# libubus/ubus-sample02/files/etc/init.d/ubus-sample02.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

start_service() {
	echo 'ubus-sample02 start'
	procd_open_instance 'ubus-sample02'
	procd_set_param command /usr/bin/ubus-sample02
	procd_close_instance
}

