# libubus/ubus-sample/files/etc/init.d/ubus-sample.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

start_service() {
	echo 'ubus-sample start'
	procd_open_instance 'ubus-sample'
	procd_set_param command /usr/bin/ubus-sample
	procd_close_instance
}

